# usr/lib/systemd/system/bashtail.service
[Unit]
Description=Periodically send logs to Loki

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bashtail
